# ssh_tunnel/etc/services.d/ssh-tunnel/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Run in s6
echo "Starting SSH tunnel service..."
exec /run.sh
